# deploy/ssh-ca-agent.service
[Unit]
Description=SSH CA Trust Agent
After=network.target

[Service]
ExecStart=/usr/local/bin/ssh-ca-agent -url http://YOUR_CA_SERVER:8080
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
